# deploy/sporttia-zero/sporttia-zero-api.service
[Unit]
Description=Sporttia ZERO API Server
After=network.target

[Service]
Type=simple
User=sporttia
Group=sporttia
WorkingDirectory=/opt/sporttia-zero/api
ExecStart=/opt/sporttia-zero/api/node_modules/.bin/tsx src/index.ts
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sporttia-zero-api

# Environment
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=/opt/sporttia-zero/api/.env

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/sporttia-zero/api

[Install]
WantedBy=multi-user.target
